# systemd/user/compton.service
[Unit]
Description=Composite manager

[Service]
ExecStart=/bin/compton
Restart=always

[Install]
WantedBy=wm.target
