# systemd/user/daily-backup.timer
[Unit]
Description=Daily backup timer

[Timer]
OnCalendar=daily
Unit=daily-backup.service

[Install]
WantedBy=timers.target
