# systemd/user/hourly-backup.timer
[Unit]
Description=Hourly backup timer

[Timer]
OnCalendar=hourly
Unit=hourly-backup.service

[Install]
WantedBy=wm.target
