# systemd/user/mbsync.timer
[Unit]
Description=Synchronize IMAP mails

[Timer]
OnBootSec=1m
OnUnitActiveSec=7m

[Install]
WantedBy=timers.target
